# docker/crontab
SHELL=/bin/bash
PATH=/opt/ai-counter/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
HOME=/home/counter
# Replaced at container start by docker/setup-timezone.sh (from TZ env or config schedule.timezone)
TZ=UTC

# Mon-Fri 06:30 — interpreted in TZ above (set schedule.timezone in config.yaml)
30 6 * * 1-5 counter /opt/ai-counter/docker/run-daily.sh >> /home/counter/ai-counter/logs/cron.log 2>&1
